# scripts/vm-docker/sync-xml-to-gitea.service
[Unit]
Description=Snapshot /opt/mcp-xml/<project>/src to Gitea admin/<Project>-src
Documentation=https://<vm-docker-host>:3000/admin/Orchestrator/src/branch/main/scripts/vm-docker/

[Service]
Type=oneshot
ExecStart=/opt/mcp-xml/sync-xml-to-gitea.sh
StandardOutput=append:/var/log/mcp-xml-sync.log
StandardError=append:/var/log/mcp-xml-sync.log
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
